# files/supervisord-sensu-client.conf
[program:sensu-client]
command = /sensu-client.sh
startretries = 10
startsecs = 35
